# etc/crontab
0 * * * * root python3 /app/scripts/update.py >/proc/1/fd/1 2>/proc/1/fd/2
0 4 * * * root python3 /app/scripts/cleanup.py >/proc/1/fd/1 2>/proc/1/fd/2
